# etc/udev/rules.d/50-powersave.rules
# Suspend when battery is at 5%
SUBSYSTEM=="power_supply", ATTR{status}=="Discharging", ATTR{capacity}=="5", RUN+="/usr/bin/systemctl suspend"
